# docroot/sites/all/modules/valkyrie/vagrant/ansible/roles/getvalkyrie.skynet/templates/skynet-queue.conf.j2
# This file is generated by Ansible. Manual changes will be overwritten!
[program:skynet-queue]
command={{ aegir_root }}/.drush/skynet/skynet.py queued --config_file={{ aegir_root }}/config/skynet.conf
numprocs=1
numprocs_start=0
priority=999
autostart=true
autorestart=unexpected
startsecs=1
startretries=3
exitcodes=0,2
stopsignal=TERM
stopwaitsecs=10
killasgroup=false
user=aegir
redirect_stderr=false
stdout_logfile=/var/log/supervisor/skynet-queue.out
stdout_logfile_maxbytes=250MB
stdout_logfile_backups=10
stderr_logfile=/var/log/supervisor/skynet-queue.err
stderr_logfile_maxbytes=250MB
stderr_logfile_backups=10
environment=HOME='{{ aegir_root }}',USER='{{ aegir_user }}'

