# docroot/sites/all/modules/valkyrie/vagrant/ansible/roles/getvalkyrie.valkyrie/templates/avahi-aliases.conf.j2
# This file is generated by Ansible. Manual changes will be overwritten!
[program:avahi-aliases]
command={{ aegir_root }}/.drush/provision_avahi/publish.py --directory={{ aegir_root }}/config/avahi-aliases --debug
numprocs=1
numprocs_start=0
priority=999
autostart=true
autorestart=unexpected
startsecs=1
startretries=3
exitcodes=0,2
stopsignal=TERM
stopwaitsecs=10
killasgroup=false
user=root
redirect_stderr=false
stdout_logfile=/var/log/supervisor/avahi-aliases.out
stdout_logfile_maxbytes=250MB
stdout_logfile_backups=10
stderr_logfile=/var/log/supervisor/avahi-aliases.err
stderr_logfile_maxbytes=250MB
stderr_logfile_backups=10

